# deploy/pev-contract-index-refresh.service
[Unit]
Description=pev, contract_index aggregate refresh (last-seen / first-seen per address)
Documentation=https://github.com/Silk-Nodes/pev
After=network-online.target postgresql.service
Wants=network-online.target

RefuseManualStart=no

[Service]
# One-shot: starts, runs the script, exits. Fired by the timer below.
Type=oneshot

User=deploy
Group=deploy
WorkingDirectory=/home/deploy/pev

ExecStart=/home/deploy/pev/deploy/run.sh db:refresh-contract-index

# Hard cap. The unnest+GROUP BY across tx_executions can be 30-90s
# depending on table size; we give it 5 minutes of headroom and rely
# on the timer's OnUnitInactiveSec=15min to keep this from overlapping
# its own next tick.
TimeoutStartSec=300s

StandardOutput=journal
StandardError=journal
SyslogIdentifier=pev-contract-index-refresh

# Security hardening (matches the other pev units)
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/home/deploy/pev
